# crontab
* * * * * /usr/bin/docker exec learntrack-app-1 php /learntrack/LearnTrack/artisan schedule:run >> /dev/null 2>&1
